# example/qbit-autodelete.timer
[Unit]
Description=Verifica a politica do qbit-autodelete a cada hora

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=5m
RandomizedDelaySec=10m
Unit=qbit-autodelete.service

[Install]
WantedBy=timers.target
